# deploy/chan.service
# /etc/systemd/system/chan.service
#
# Systemd service for the Chan imageboard.
#
# Install:
#   sudo cp chan.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable chan
#   sudo systemctl start chan

[Unit]
Description=Chan Imageboard
Documentation=https://github.com/yourname/chan
After=network.target
Wants=network.target

[Service]
# Run as a dedicated unprivileged user for security isolation
Type=simple
User=chan
Group=chan

# Binary and working directory
ExecStart=/usr/local/bin/chan
WorkingDirectory=/var/lib/chan

# Restart policy: restart on failure after 5 seconds
Restart=on-failure
RestartSec=5s

# ── Environment ────────────────────────────────────────────────────────────────
# Override in /etc/systemd/system/chan.service.d/override.conf
Environment=CHAN_BIND=0.0.0.0:8080
Environment=CHAN_DB=/var/lib/chan/chan.db
Environment=CHAN_UPLOADS=/var/lib/chan/uploads

# IMPORTANT: Change this secret before first run!
# Generate a good one with: openssl rand -hex 32
Environment=CHAN_COOKIE_SECRET=CHANGE_ME_TO_A_RANDOM_64_CHAR_HEX_STRING

# Optional tuning
Environment=CHAN_MAX_FILE_SIZE=8388608
Environment=CHAN_THUMB_SIZE=250
Environment=CHAN_RATE_POSTS=10
Environment=CHAN_RATE_WINDOW=60
Environment=CHAN_BEHIND_PROXY=false
Environment=RUST_LOG=chan=info

# ── Hardening ──────────────────────────────────────────────────────────────────
# Prevent writing anywhere except the data directory
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/chan

# No new privileges (prevents suid escalation)
NoNewPrivileges=true

# Private /tmp
PrivateTmp=true

# Restrict available system calls to a safe subset
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Restrict network namespaces
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true

# Memory execution prevention
MemoryDenyWriteExecute=true

# Deny access to kernel module loading
ProtectKernelModules=true
ProtectKernelTunables=true

# Limit to what we actually need
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits: protect against runaway process on Pi
LimitNOFILE=4096
LimitNPROC=128

# Graceful shutdown timeout
TimeoutStopSec=15s

# Log to journald (no separate log file → less SD card wear)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chan

[Install]
WantedBy=multi-user.target
